# caretaker/digiquarium-caretaker.service
[Unit]
Description=Digiquarium Autonomous Caretaker
After=network.target

[Service]
Type=simple
User=ijneb
WorkingDirectory=${DIGIQUARIUM_HOME:-/home/ijneb/digiquarium}/caretaker
Environment=DIGIQUARIUM_DIR=${DIGIQUARIUM_HOME:-/home/ijneb/digiquarium}
ExecStart=/usr/bin/python3 ${DIGIQUARIUM_HOME:-/home/ijneb/digiquarium}/caretaker/caretaker.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
