# daemons/webmaster/digiquarium-broadcaster.service
[Unit]
Description=Digiquarium Live Feed Broadcaster
After=network.target docker.service

[Service]
Type=simple
User=ijneb
WorkingDirectory=${DIGIQUARIUM_HOME:-/home/ijneb/digiquarium}/daemons/webmaster
Environment=DIGIQUARIUM_HOME=/home/ijneb/digiquarium
ExecStart=/usr/bin/python3 ${DIGIQUARIUM_HOME:-/home/ijneb/digiquarium}/daemons/webmaster/broadcaster_continuous.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
